# test/docker-conf/phpfpm/supervisord_workers.conf


[program:phpfpm]
command=/usr/sbin/php-fpm -F -O
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm.log
stopsignal=INT
priority=1

[program:resque-server]
command=/usr/bin/php /app/test/testapp/console.php resque:server
autostart=true
autorestart=true
user=userphp
umask=022
directory=/tmp
stopsignal=QUIT
redirect_stderr=true
stdout_logfile=/var/log/resque-server.log
priority=2
